# deploy/systemd/sentinel-stakeholder.service.example
# Sentinel stakeholder report systemd service template.
# Template only. Do not install, enable or start until Raspberry Pi deployment controls are agreed.
# Copy to /etc/systemd/system/sentinel-stakeholder.service only during a controlled deployment.

[Unit]
Description=Sentinel weekly stakeholder progress report
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/sentinel/apps/seo-ops
EnvironmentFile=/srv/sentinel/apps/seo-ops/.env
Environment=PLATFORM_TENANT=erp-experts
ExecStart=/usr/bin/npm run platform:cadence -- --stakeholder-only
User=matthew
Group=matthew
