# deploy/systemd/sentinel-api.service.example
# Sentinel API systemd service template.
# Template only. Do not install or enable until auth and Raspberry Pi deployment controls are agreed.
# Copy to /etc/systemd/system/sentinel-api.service only during a controlled deployment.

[Unit]
Description=Sentinel local API service
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/sentinel/apps/seo-ops
EnvironmentFile=/srv/sentinel/apps/seo-ops/.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=SENTINEL_API_HOST=127.0.0.1
Environment=SENTINEL_API_PORT=4317
Environment=PLATFORM_TENANT=erp-experts
ExecStart=/usr/local/bin/npm run platform:api:serve
Restart=on-failure
RestartSec=5
User=matthew
Group=matthew

[Install]
WantedBy=multi-user.target
